# ansible/dropbear/dropbear.service
[Unit]
Description=Dropbear SSH Server Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/dropbear -R

[Install]
WantedBy=multi-user.target